# deploy/media-pi-daemon.service
# Media-Pi daemon — upload pipeline (scheduler + MQTT to follow).
# Install on the Pi:
#   sudo cp deploy/media-pi-daemon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now media-pi-daemon
# Logs: journalctl -u media-pi-daemon -f

[Unit]
Description=Media-Pi upload pipeline daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/media-pi
EnvironmentFile=/opt/media-pi/.env
# daemon/tui live under src/ and there's no editable install on the Pi — put src
# on the path so `python -m daemon` (and its `import tui.*`) resolve. Mirrors the
# sys.path insert the scripts/ use.
Environment=PYTHONPATH=/opt/media-pi/src
# Apply the schema and seed status_mapping before start; the daemon hard-exits if
# the required statuses are missing. Idempotent (INSERT OR IGNORE), safe every boot.
ExecStartPre=/usr/bin/python3 scripts/seed_constants.py
ExecStart=/usr/bin/python3 -m daemon
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
